# .config/systemd/user/mbsync@.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnCalendar=*:0/10
Unit=mbsync@%i.service

[Install]
WantedBy=timers.target
